# .config/supervisor/app.conf
[program:nginx]
command = sudo nginx -g 'daemon off;'