# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = sudo uwsgi --ini /srv/project/.config/uwsgi.ini